# ccso.socket
# /etc/systemd/system/ccso.socket
[Unit]
Description=CCSO Nameserver

[Socket]
ListenStream=105
Accept=yes

[Install]
WantedBy=sockets.target